# claudio.service
[Unit]
Description=Claudio Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-bot
ExecStart=/usr/local/bin/gunicorn --workers 1 --bind 0.0.0.0:5000 --timeout 120 --preload app:app
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
